# chat_bot/deploy/papp-chat.service
[Unit]
Description=PAPP EA Chat Assistant (FastAPI)
After=network.target postgresql.service opencode-serve.service
Wants=postgresql.service opencode-serve.service

[Service]
Type=simple
User=pietro_giacobazzi
WorkingDirectory=/home/pietro_giacobazzi/Desktop/PAPP_EA/chat_bot
EnvironmentFile=/home/pietro_giacobazzi/Desktop/PAPP_EA/chat_bot/.env
# La chat invoca il binario `opencode`: deve essere nel PATH del servizio.
Environment=PATH=/home/pietro_giacobazzi/.opencode/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/python3 -m uvicorn app:app --host 0.0.0.0 --port 8090
Restart=always
RestartSec=3
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
