# chat_bot/deploy/papp-backup.timer
[Unit]
Description=Esegue il backup DB PAPP EA ogni giorno

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true

[Install]
WantedBy=timers.target
